# 1/06/systemd/system/httpd.service.d/override.conf
[Service]
PIDFile=/var/run/httpd/%i.pid
EnvironmentFile=/etc/sysconfig/httpd-%i
