# 1/06/systemd/system/httpd@.service.d/override.conf
[Service]
EnvironmentFile=/etc/sysconfig/httpd-%i